# config/services/n8n.service
[Unit]
Description=apinasaari-slackbots n8n service
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=n8n
Group=n8n
WorkingDirectory=/opt/apinasaari-slackbots
Environment="PODMAN_USERNS=keep-id"
Environment="XDG_RUNTIME_DIR=/run/user/1001"
ExecStartPre=/usr/bin/podman-compose pull
ExecStart=/usr/bin/podman-compose up --no-deps --remove-orphans
ExecStop=/usr/bin/podman-compose down
KillMode=mixed
KillSignal=SIGINT
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=60

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/apinasaari-slackbots /home/n8n

[Install]
WantedBy=multi-user.target
